# dot-config/systemd/user/himalaya-watch-gmail.service
# ~/.config/systemd/user/himalaya-watch.service

[Unit]
Description=Himalaya watcher for gmail
After=network.target

[Service]
ExecStart=/home/djha-skin/.local/bin/him -a gmail watch
Environment=RUST_LOG=debug
Restart=always
RestartSec=31

[Install]
WantedBy=default.target