# dot-config/systemd/user/onedrive-sync.service
[Unit]
Description=OneDrive synchronization service

[Service]
Type=oneshot
ExecStart=/home/skin/.local/bin/onedrive-sync