# dot-config/systemd/user/sync-mail.timer
[Unit]
Description=Mailbox synchronization timer

[Timer]
OnBootSec=2m
OnUnitActiveSec=20m
Unit=sync-mail.service

[Install]
WantedBy=timers.target